# deploy/systemd/privacy-calling-api.service
[Unit]
Description=Privacy Calling API
After=network.target asterisk.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/fyp/PrivacyCalling
EnvironmentFile=/home/ubuntu/fyp/PrivacyCalling/deploy/env/privacy-calling.env
ExecStart=/usr/bin/node /home/ubuntu/fyp/PrivacyCalling/src/index.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
